# files/ftp-mirror.service
[Unit]
Description="FTP Mirror"

[Service]
ExecStart=/usr/local/bin/ftp-mirror /usr/local/etc/ftp-mirror-config.yaml
